# system/zram-writeback.timer
[Unit]
Description=Timer for zram-writeback
After=multi-user.target

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
